# deploy/ae-template-site.service
[Unit]
Description=AE Template Library
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/ae-template-site
Environment=NODE_ENV=production
EnvironmentFile=-/opt/ae-template-site/.env.local
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
User=ae-template-site
Group=ae-template-site
UMask=0077

[Install]
WantedBy=multi-user.target
